# BlowLeafAnsible/roles/deploy-blowleaf/templates/back.service.j2
[Unit]
Description=BlowLeaf back service

[Service]
User={{user}}
Group={{group}}
ExecStart=java -jar {{back.service.path}}/{{back.service.bin}}

[Install]
WantedBy=multi-user.target